# server2/nebula
#!/sbin/openrc-run

name="Nebula"
description="Nebula overlay networking tool"
command="/usr/local/bin/nebula"
command_args="-config /etc/nebula/config.yml"
command_background="yes"
pidfile="/run/$RC_SVCNAME.pid"

depend() {
    need net
    after firewall
}

start_pre() {
    checkpath -f -m 0600 -o nebula:nebula /etc/nebula/config.yml
}

reload() {
    ebegin "Reloading $name"
    start-stop-daemon --signal HUP --pidfile "$pidfile"
    eend $?
}